# firmware/rootfs/57.10-25140/etc/inittab
# autogenerated by gen_inittab.py for ARCH fenway; DO NOT EDIT
::sysinit:/etc/Configure
null::respawn:/wifi/netstartd
null::respawn:/etc/runmdns
null::respawn:/etc/runanacapa
null::respawn:/etc/runchrony
null::respawn:/etc/runsddp
/dev/ttyS0::askfirst:/bin/sh
::ctrlaltdel:/sbin/reboot
::shutdown:/etc/init.d/rcK
::restart:/sbin/init
