# firmware/rootfs/16.6-00002-diag/etc/inittab
::sysinit:/etc/Configure
null::respawn:/etc/runanacapa
/dev/ttyS0::askfirst:/bin/sh
::ctrlaltdel:/sbin/reboot
::shutdown:/etc/init.d/rcK
::restart:/sbin/init
