# firmware/rootfs/28.1-83040/etc/inittab
::sysinit:/etc/Configure
null::respawn:/wifi/netstartd
null::respawn:/etc/runanacapa
::askfirst:/bin/sh
::ctrlaltdel:/sbin/reboot
::shutdown:/etc/init.d/rcK
::restart:/sbin/init
